# userspace/inittab
/bin/larpbox:getty:/dev/tty0
/bin/larpbox:getty:/dev/tty1
